# installer/ourbox/rootfs/etc/update-motd.d/10-ourbox-status
#!/bin/bash
STATUS_FILE="/run/ourbox/status"
if [[ -f "$STATUS_FILE" ]]; then
  echo ""
  cat "$STATUS_FILE"
  echo ""
else
  echo ""
  echo "OurBox Woodbox: status service starting..."
  echo "Run \"woodbox status\" to check, or wait a moment and reconnect."
  echo ""
fi
